# systemd-user/pikachu-n8n-backup.timer
[Unit]
Description=Daily backup timer (Pikachu n8n)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
