# systemd-user/pikachu-n8n-backup.service
[Unit]
Description=Pikachu n8n encrypted backup to pCloud

[Service]
Type=oneshot
WorkingDirectory=/code/pikachu
ExecStart=/code/pikachu/scripts/backup.sh