# systemd-user/pikachu-n8n.service
[Unit]
Description=Pikachu n8n stack (podman-compose)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/code/pikachu
ExecStart=/code/pikachu/scripts/up.sh
ExecStop=/code/pikachu/scripts/down.sh
TimeoutStartSec=300
TimeoutStopSec=300

[Install]
WantedBy=default.target
